# dot_local/bin/dev.kanata.wake-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Wake WATCHDOG: closes the gap that the sleepwatcher fast-path (dev.kanata.wake)
         misses. sleepwatcher's -w fires only on FULL user wakes; the kanata input grab
         can also die silently on DarkWake / scheduled-alarm wakes (e.g. a calaccessd
         alarm). This polls every 15s and re-grabs whenever pmset's most-recent wake is
         newer than the kanata process start. launchd defers StartInterval while asleep,
         so it costs ~nothing on battery and never wakes a sleeping machine.
         Logic lives in kanata-wake-watchdog.sh; it execs kanata-wake-restart.sh, which
         carries the self-expiring lock that dedupes against the sleepwatcher hook. -->
    <key>Label</key><string>dev.kanata.wake-watchdog</string>
    <key>UserName</key><string>root</string>
    <key>ProgramArguments</key>
    <array><string>/Users/maximegardoni/.local/bin/kanata-wake-watchdog.sh</string></array>
    <key>StartInterval</key><integer>15</integer>
    <key>RunAtLoad</key><true/>
    <key>StandardOutPath</key><string>/var/log/kanata-wake.log</string>
    <key>StandardErrorPath</key><string>/var/log/kanata-wake.log</string>
</dict>
</plist>
